# custom/salt-gitfs-http/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:uwsgi]
priority=100
command=/usr/sbin/uwsgi --yaml /etc/uwsgi/gitfs.yml

[program:nginx]
priority=200
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
